# miniforge_clean/pkgs/conda-25.3.1-py312h7900ff3_1/etc/profile.d/conda.csh
setenv CONDA_EXE "/home/conda/feedstock_root/build_artifacts/conda_1744220005005/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_place/bin/conda";
setenv _CONDA_ROOT "/home/conda/feedstock_root/build_artifacts/conda_1744220005005/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_place";
setenv _CONDA_EXE "/home/conda/feedstock_root/build_artifacts/conda_1744220005005/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_place/bin/conda";
setenv CONDA_PYTHON_EXE "/home/conda/feedstock_root/build_artifacts/conda_1744220005005/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_place/bin/python";
source "/home/conda/feedstock_root/build_artifacts/conda_1744220005005/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_place/lib/python3.12/site-packages/conda/shell/etc/profile.d/conda.csh";
