# miniforge_clean/etc/profile.d/conda.csh
setenv CONDA_EXE "/home/luiz/Área de trabalho/EMPRESAS/UNISAGRADO/PESQUISAS/VASSOURA DE BRUXA/miniforge_clean/bin/conda";
setenv _CONDA_ROOT "/home/luiz/Área de trabalho/EMPRESAS/UNISAGRADO/PESQUISAS/VASSOURA DE BRUXA/miniforge_clean";
setenv _CONDA_EXE "/home/luiz/Área de trabalho/EMPRESAS/UNISAGRADO/PESQUISAS/VASSOURA DE BRUXA/miniforge_clean/bin/conda";
setenv CONDA_PYTHON_EXE "/home/luiz/Área de trabalho/EMPRESAS/UNISAGRADO/PESQUISAS/VASSOURA DE BRUXA/miniforge_clean/bin/python";
source "/home/luiz/Área de trabalho/EMPRESAS/UNISAGRADO/PESQUISAS/VASSOURA DE BRUXA/miniforge_clean/lib/python3.12/site-packages/conda/shell/etc/profile.d/conda.csh";
